# config/autostart/disable-x11-screensaver.desktop
[Desktop Entry]
Type=Application
Name=Disable X11 Screen Saver
Comment=Disable the X server built-in screen saver so xfce4-screensaver is the sole authority for idle locking
Exec=/usr/local/bin/disable-x11-screensaver.sh
StartupNotify=false
NoDisplay=true
